# app-admin/amazon-cloudwatch-agent/files/amazon-cloudwatch-agent.initd
#!/sbin/openrc-run
# Copyright 2026 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Amazon CloudWatch Agent"
command="/usr/bin/amazon-cloudwatch-agent"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}
